# .config/local/supervisor_app.conf
[program:uwsgi]
command=uwsgi --ini /srv/project/.config/local/uwsgi.ini

[program:front]
command=ng serve
directory=/srv/front

[program:nginx]
command=nginx